# backend/crates/kino/rpm/kino.service
[Unit]
Description=Kino — single-binary media automation and streaming server
Documentation=https://kinostack.app
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kino
Group=kino
# Opt the binary into "exit-after-restore" so a successful POST to
# /api/v1/backups/{id}/restore exits with EX_TEMPFAIL and systemd
# restarts the process against the freshly-restored database
# (Restart=on-failure below catches the non-zero exit). No user
# action required.
# `--data-path` and `--no-open-browser` are TOP-LEVEL flags on
# `kino`, not on the `serve` subcommand — passing them after `serve`
# fails clap parsing with status=2/INVALIDARGUMENT. Setting them via
# the documented env vars sidesteps the order-dependence entirely.
Environment=KINO_RESTART_AFTER_RESTORE=1
Environment=KINO_DATA_PATH=/var/lib/kino
Environment=KINO_NO_OPEN_BROWSER=1
Environment=KINO_INSTALL_KIND=linux-systemd
# CWD pinned to a writable kino-owned dir — without this, relative
# I/O paths (notably librqbit's first-piece open) try to write
# under whatever CWD systemd inherited (often /home/<user>),
# which is empty inside the namespace under ProtectHome=true.
WorkingDirectory=/var/lib/kino
# Port 80 is the schema default for fresh installs.
# CAP_NET_BIND_SERVICE below lets the unprivileged kino user bind
# it. Don't set KINO_PORT here — it's consulted only at first-run
# schema insert; setting it permanently would mask user changes
# from Settings → General → Port.
ExecStart=/usr/bin/kino serve
Restart=on-failure
RestartSec=5
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/kino /etc/kino
PrivateTmp=true
RuntimeDirectory=kino
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
